# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/opt/liferay/tomcat/bin/catalina.sh run
numprocs=1
autostart=true
autorestart=true